# server/etc/systemd/system/restart-systemd-networkd.path
[Path]
PathChanged=/etc/systemd/network/999-wg0.netdev.d

[Install]
WantedBy=network.target
